# 0x0-prune.timer
[Unit]
Description=Prune 0x0 files

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
